# etc/rc.local
#!/bin/sh -e

# Fix ld.so.preload from raspi-copies-and-fills
mv /etc/ld.so.preload.disabled /etc/ld.so.preload

echo "#!/bin/sh -e
exit 0
" > /etc/rc.local

exit 0
