# etc/rc.local
#!/bin/sh -e


echo "#!/bin/sh -e
exit 0
" > /etc/rc.local

exit 0
